# systemd/monitoring-agent.service
[Unit]
Description=HostMonitor Agent
After=network.target

[Service]
Type=simple
# Запуск от root для доступа к системным метрикам
# User=root убран из-за ошибки 217/USER в некоторых конфигурациях systemd
WorkingDirectory=/opt/monitoring/agent
Environment="PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/opt/monitoring/.venv/bin"
# Переменные будут загружены из node.conf через load_node_conf() в main.py
ExecStart=/opt/monitoring/.venv/bin/python3 /opt/monitoring/agent/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
# Увеличиваем таймаут для запуска
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target

